# debian/aptik-battery-monitor.postinst
#!/bin/sh

#DEBHELPER#

set -e

#create log directory
mkdir -p /var/log/aptik-battery-monitor  
chmod a+rwx /var/log/aptik-battery-monitor

#start service
/etc/init.d/aptik-battery-monitor start
sudo update-rc.d aptik-battery-monitor defaults
